# packer-builder-arm/hack-a-sat/comm/systemd/services/set-ethernet-speed.service
[Unit]
Description=Set Ethernet Speed to 100BaseTX on Boot

[Service]
ExecStart=/usr/local/sbin/set-ethernet-speed.sh

[Install]
WantedBy=network-online.target
